# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:stats]
command=/usr/bin/node stats.js
directory=/opt/docker-stats
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:httpd]
command=/usr/bin/node httpd.js
directory=/opt/docker-stats
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
